# psscan.upstart
# psscan ps logging daemon
#

description "psscan - ps logging daemon"

start on runlevel [345]
stop on runlevel [016]

respawn
exec /usr/sbin/psscan 
